# templates/simple_watchdog.timer.j2
[Unit]
Description={{ ya_requestor_name }} simple watchdog

[Timer]
OnCalendar={{ ya_requestor_simple_watchdog_when }}

[Install]
WantedBy=timers.target
